# build/build/workspace/etc/systemd/system/aurora-mcp.service
[Unit]
Description=Aurora OS MCP Nervous System
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aurora/mcp
ExecStart=/usr/bin/python3 system/mcp_host.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
